# plugins/mc-board/agent-runner/com.miniclaw.board-agent-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.miniclaw.board-agent-runner</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/augmentedmike/am/miniclaw/plugins/mc-board/agent-runner/runner.mjs</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/augmentedmike</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>OPENCLAW_STATE_DIR</key>
    <string>/Users/augmentedmike/am</string>
    <key>OPENCLAW_BIN</key>
    <string>/opt/homebrew/bin/openclaw</string>
    <key>CLAUDE_BIN</key>
    <string>/opt/homebrew/bin/claude</string>
    <key>BOARD_DB_PATH</key>
    <string>/Users/augmentedmike/am/user/augmentedmike_bot/brain/board.db</string>
    <key>AGENT_RUNNER_POLL_MS</key>
    <string>5000</string>
    <key>AGENT_RUNNER_MAX_CONCURRENT</key>
    <string>3</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/augmentedmike/am/logs/agent-runner.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/augmentedmike/am/logs/agent-runner.stderr.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/augmentedmike/am/miniclaw/plugins/mc-board/agent-runner</string>
</dict>
</plist>
