# com.miniclaw.board-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.miniclaw.board-web</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/michaeloneal/.nvm/versions/node/v24.9.0/bin/node</string>
    <string>/Users/michaeloneal/.openclaw/miniclaw/web/node_modules/.bin/next</string>
    <string>start</string>
    <string>-p</string>
    <string>4220</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/michaeloneal/.openclaw/miniclaw/web</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>StandardOutPath</key>
  <string>/Users/michaeloneal/.openclaw/logs/miniclaw-board-web.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/michaeloneal/.openclaw/logs/miniclaw-board-web.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/michaeloneal</string>
    <key>PATH</key>
    <string>/Users/michaeloneal/.local/bin:/Users/michaeloneal/.nvm/versions/node/v24.9.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>OPENCLAW_STATE_DIR</key>
    <string>/Users/michaeloneal/.openclaw</string>
  </dict>
</dict>
</plist>
